# debian/postinst
#! /bin/sh

set -e

add_skolab_system_users() {
	if ! getent passwd skolab >/dev/null; then
		adduser --group --system --no-create-home --gecos "Skolab Groupware system account" \
			--home /var/lib/skolab skolab
	fi
}

adduser_postfix_in_sasl() {
	# Add user postfix to group sasl, if such group exists
        adduser postfix sasl || true
}

adduser_cyrus_in_skolab() {
	# Add user cyrus to group skolab, if such group exists
	# This is needed for PKI
        adduser cyrus skolab || true
}

fixperms() {
	chown skolab:skolab /var/lib/skolab
	chmod 2775 /var/lib/skolab
	
	chmod 750 /var/lib/skolab/freebusy-cache/
	chown www-data:root /var/lib/skolab/freebusy-cache/
}

configure_apache() {
	#a2enmod and a2ensite are in same package apache*-common
	if command -v a2enmod >/dev/null 2>&1; then
		for mod in rewrite ssl ldap; do
			a2enmod $mod >/dev/null
		done

		if [ -e /etc/apache2/mods-available/php7.load ]; then
			a2enmod php7 >/dev/null
		fi

		a2ensite skolab >/dev/null
	fi

	if command -v invoke-rc.d >/dev/null 2>&1; then
		invoke-rc.d apache2 force-reload || :
	else
		[ -x /etc/init.d/apache2 ] && /etc/init.d/apache2 force-reload || :
	fi
}

if [ "$1" = configure ]; then
	add_skolab_system_users
	adduser_postfix_in_sasl
	adduser_cyrus_in_skolab
	fixperms
	configure_apache
fi

#DEBHELPER#
